# config/supervisor.conf
[program:api1]
command = /bin/bash -c "python /app/displacy/app.py"

[program:api2]
command = /bin/bash -c "python /app/matcher/app.py"