# config/supervisor.conf
[program:api]
command = /usr/bin/make --directory /app start
